# skills/h-v3-bot/h_v3.service
[Unit]
Description=H_V3 Dolphin AI Strategy Engine (MCP Protocol)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/h_v3
ExecStart=/usr/bin/python3.11 /root/h_v3/h_v3_bot.py
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/tmp/h_v3.log
StandardError=append:/tmp/h_v3.log

[Install]
WantedBy=multi-user.target
